# deploy/mem-manager.service
[Unit]
Description=mem-manager — Claude memory MCP server
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=masud
WorkingDirectory=/home/masud/mem-manager
ExecStart=/usr/bin/node /home/masud/mem-manager/scripts/mcp-server.js
Restart=always
RestartSec=5
EnvironmentFile=/home/masud/mem-manager/.env
StandardInput=null
StandardOutput=journal
StandardError=journal
SyslogIdentifier=mem-manager
KillMode=process

[Install]
WantedBy=multi-user.target
